# scripts/com.cave.upstream-watchd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  cave-upstream-watchd launchd plist (macOS).

  Install:
    cp scripts/com.cave.upstream-watchd.plist ~/Library/LaunchAgents/
    launchctl load -w ~/Library/LaunchAgents/com.cave.upstream-watchd.plist

  Uninstall:
    launchctl unload -w ~/Library/LaunchAgents/com.cave.upstream-watchd.plist
    rm ~/Library/LaunchAgents/com.cave.upstream-watchd.plist

  Logs:
    ~/Library/Application Support/cave-runtime/watchd/watchd.log
    ~/Library/Application Support/cave-runtime/watchd/watchd.err

  Environment:
    The GitHub PAT is read by cave-upstream-watchd via the macOS
    keychain (service `cave-upstream-watchd`, account `$USER`):

      security add-generic-password -U \
        -s cave-upstream-watchd \
        -a "$USER" \
        -w "<your-PAT>"

    If the keychain item is missing the daemon falls back to
    `$GITHUB_TOKEN` (DEPRECATED — emits a warn log), then to
    anonymous 60 req/h. Do NOT commit a token in
    `EnvironmentVariables`.

  Cadence:
    StartInterval = 300 seconds = 5 minutes. Matches Burak's
    self-healing/self-improving Charter: any new upstream release is
    detected and surfaced as a GAP_OPENED event within five minutes
    of being published.

  Linux equivalent:
    A systemd timer + service should be added under
    `deploy/systemd/cave-upstream-watchd.{service,timer}` when the
    daemon is ready for the Linux production path. The TOML config
    file shape + binary CLI are identical.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.cave.upstream-watchd</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/cave-upstream-watchd</string>
    <string>poll</string>
  </array>

  <!-- 5-minute poll cadence per Burak's spec. -->
  <key>StartInterval</key>
  <integer>300</integer>

  <key>RunAtLoad</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/gnomish/Library/Application Support/cave-runtime/watchd/watchd.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/gnomish/Library/Application Support/cave-runtime/watchd/watchd.err</string>

  <!--
    No EnvironmentVariables{GITHUB_TOKEN} block — keychain-first PAT
    resolution is the supported path. See the header comment for the
    `security add-generic-password` command. PATH/HOME are picked up
    from the launchd user-agent context by default.
  -->

  <!-- Lower priority so the daemon doesn't compete with the runtime. -->
  <key>LowPriorityIO</key>
  <true/>
  <key>Nice</key>
  <integer>5</integer>
</dict>
</plist>
